# Ubuntu/IntelligentHub/workspaceone-intelligent-hub-amd64-21.01.0.1/rootfs/opt/Workspace-ONE-Intelligent-Hub/script/ws1Hubd
#!/sbin/openrc-run

#
# Copyright (C) 2019 VMware, Inc. All rights reserved.
# This product is protected by copyright and intellectual property laws in
# the United States and other countries as well as by international treaties.
# VMware products may be covered by one or more patents listed at
# https://www.vmware.com/go/patents.
#

### BEGIN INIT INFO
# Provides:             ws1Hubd
# Required-Start:       $all
# Required-Stop:
# Default-Start:        3 4 5
# Default-Stop:
# Short-Description:    Workspace ONE Intelligent Hub agent restart
### END INIT INFO

name="ws1Hubd"
description="Open-RC service for Workspace ONE Intellient Hub to start the agent on system reboot"
AGENT_NAME="ws1Hub"
AGENT_CONF=/etc/$CONF_NAME.conf
WS1HOME=$(awk -F '[<>]' '/<home>/ {print $3}' $AGENT_CONF)

start() {
  "$WS1HOME"/agent start 
}

stop() {
  "$WS1HOME"/agent stop 
}


